# docker/python/3/cron/crontab
*/1 * * * * root python3 /src/hello.py >> /var/log/cron.log
# Waiting for Shishkebabs!
