# docker/centos/supervisor/supervisord.conf
[supervisord]
nodaemon=true
logfile = /tmp/supervisord.log
pidfile = /tmp/supervisord.pid

[program:sshd]
command=/usr/sbin/sshd -D
files=/etc/ssh/sshd_config /etc/ssh/ssh_config
